# backend/nse-api.service
[Unit]
Description=NSE Trading API (FastAPI)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/nse-trading-signals/backend
ExecStart=/usr/bin/python3 -m uvicorn main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=5
EnvironmentFile=/root/nse-trading-signals/.env

[Install]
WantedBy=multi-user.target
